# purchasesbot.service
[Unit]
Description=Element Bot - Purchases Telegram Bot
After=network.target

[Service]
Type=simple
User=your_user
WorkingDirectory=/path/to/elementBot
Environment="PATH=/path/to/elementBot/venv/bin"
ExecStart=/path/to/elementBot/venv/bin/python /path/to/elementBot/beats_purchases_bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target


